# tests/fixtures/systemd/myunit.socket
[Socket]
ListenStream=/run/{{unit}}/socket
MaxConnections=4
